# nginx/console/console.service
[Unit]
Description=uWSGI instance to serve console

[Service]
ExecStart=/services/console/serve.sh

[Install]
WantedBy=multi-user.target
