# nginx/api/api.service
[Unit]
Description=uWSGI instance to serve api

[Service]
ExecStart=/services/api/serve.sh

[Install]
WantedBy=multi-user.target
